# ansible/roles/wiki/templates/wiki_backup.service.j2
[Unit]
Description=Backup wiki data

[Service]
Type=oneshot
User={{ data_service_user }}
ExecStart=/var/scripts/wiki_backup.sh

[Install]
WantedBy=multi-user.target