# ansible/roles/wiki/templates/wiki_backup.timer.j2
[Unit]
Description=Backup wiki data timer

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target